# test/rc.local
#!/bin/bash

echo "script-test: rc.local executed"

exit 0
